# config/systemd/user/ryoku-shell.service.d/qt6-fractional-scale-workaround.conf
# Workaround for Qt 6.11.0 stack overflow in libQt6Quick triggered by
# wl_surface.preferred_buffer_scale events from niri (lid-close,
# output-config-change, etc.).
#
# The recursion happens at libQt6Quick.so.6 offset 0x248e81 inside
# QQuickWindow::physicalDpiChanged emission, called from
# QtWaylandClient::QWaylandWindow::setScale(d). Disabling Qt's
# fractional-scale handling routes through a different code path that
# avoids the recursive child-window DPR propagation walker.
#
# Coredump evidence: 10 quickshell SIGSEGVs in 82 minutes on 2026-05-03
# while the bug was active; zero coredumps after this drop-in landed.
#
# To revert: rm this file, systemctl --user daemon-reload,
#            systemctl --user restart ryoku-shell.service.

[Service]
Environment=QT_WAYLAND_DISABLE_FRACTIONAL_SCALE=1
