# config/systemd/user/ryoku-audio-restore-mixers.service
[Unit]
Description=Ryoku ALSA hardware mixer unmute (self-heal silent-speakers state)
After=pipewire.service

[Service]
Type=oneshot
ExecStart=%h/.local/share/ryoku/bin/ryoku-audio-restore-mixers
LogLevelMax=warning

[Install]
WantedBy=default.target
