# config/systemd/user/ryoku-resume-listener.service
[Unit]
Description=Ryoku resume listener (logind PrepareForSleep)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/share/ryoku/bin/ryoku-resume-listener
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
