# config/systemd/user/vicinae.service
[Unit]
Description=Vicinae launcher daemon
Documentation=https://docs.vicinae.com
After=graphical-session.target
StartLimitIntervalSec=30
StartLimitBurst=3

[Service]
Type=simple
# Resolve the binary from PATH so this works whether vicinae lives in the
# vicinae-bin package (/usr/bin) or a user-local install (~/.local/bin).
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin
ExecStart=/usr/bin/env vicinae server --replace
ExecReload=/bin/kill -HUP $MAINPID
SuccessExitStatus=143
KillMode=process
KillSignal=SIGTERM
Restart=on-failure
RestartSec=2

# No [Install] section. Hyprland starts this service from the Ryoku session
# config (exec-once), the same way it starts ryoku-shell.service, so vicinae
# inherits the live Wayland session environment. graphical-session.target is
# not reliably active in Ryoku Hyprland sessions, so we do not enable against
# it.
